# files/etc/init.d/vhusbd
#!/bin/sh /etc/rc.common

START=99

start() {
   /usr/bin/vhusbd -b >/dev/null 2>&1 &
}

stop() {
    killall -9 vhusbd >/dev/null 2>&1
}
